# scripts/podman/deneb.container.in
# Deneb gateway — Podman Quadlet (rootless)
# Installed by setup-podman.sh into deneb's ~/.config/containers/systemd/
# {{DENEB_HOME}} is replaced at install time.

[Unit]
Description=Deneb gateway (rootless Podman)

[Container]
Image=deneb:local
ContainerName=deneb
UserNS=keep-id
# Keep container UID/GID aligned with the invoking user so mounted config is readable.
User=%U:%G
Volume={{DENEB_HOME}}/.deneb:/home/node/.deneb:Z
EnvironmentFile={{DENEB_HOME}}/.deneb/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color
PublishPort=18789:18789
PublishPort=18790:18790
Pull=never
Exec=node dist/index.js gateway --bind lan --port 18789

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
