# scripts/deploy/wormhole.service
# Install as a systemd --user unit (matching deneb-gateway.service):
#   cp scripts/deploy/wormhole.service ~/.config/systemd/user/
#   systemctl --user daemon-reload && systemctl --user enable --now wormhole.service
# A --user unit runs as the logged-in user already, so it carries no User= (that
# is a system-unit directive and fails here with 216/GROUP) and is wanted by
# default.target (the user manager has no multi-user.target).
[Unit]
Description=wormhole model router (Deneb) — OpenAI/Anthropic-compatible multiplexer
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
# Cloud upstream keys for wormhole-routed subscription models (e.g. ${ZAI_API_KEY}
# referenced by ~/.wormhole/config.json). Optional (leading '-'): wormhole runs
# fine without it for local-only routing.
EnvironmentFile=-/home/choiceoh/.deneb/.env
# wormhole fans one endpoint out to many model backends (local vLLM + cloud) by
# model name. With Deneb's roles pointed at it, it sits in the model hot path, so
# it restarts on failure. The binary is built by `make wormhole` into dist/.
ExecStart=/home/choiceoh/deneb/dist/wormhole --config /home/choiceoh/.wormhole/config.json
WorkingDirectory=/home/choiceoh/deneb
Restart=on-failure
RestartSec=5

# A streaming reverse proxy: it buffers at most a few in-flight request bodies
# (32 MiB cap each) and streams responses through, so its own footprint is small.
# 512M is generous headroom and a guard against a runaway.
MemoryMax=512M
MemoryHigh=384M

# Graceful shutdown.
TimeoutStopSec=10
KillSignal=SIGTERM

# Logging.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wormhole

[Install]
WantedBy=default.target
